# deploy/cliproxy-access-manager.service
[Unit]
Description=CLIProxy API Access Manager Sidecar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cliproxy-access-manager
ExecStart=/opt/cliproxy-access-manager/cliproxy-access-manager
Restart=always
RestartSec=5
Environment=CLIPROXY_MGMT_BASE_URL=http://127.0.0.1:8317
Environment=APIM_HTTP_ADDR=127.0.0.1:8390
Environment=APIM_DATA_DIR=/var/lib/cliproxy-access-manager
Environment=APIM_DB_PATH=/var/lib/cliproxy-access-manager/apim.db
Environment=APIM_KEY_SYNC_INTERVAL=30s
Environment=APIM_USAGE_SYNC_INTERVAL=2m
Environment=APIM_RECOVERY_CHECK_INTERVAL=60s
Environment=APIM_UPDATE_CHECK_ENABLED=true
Environment=APIM_UPDATE_CHECK_TIME=04:00
Environment=APIM_MANAGEMENT_LATEST_VERSION_URL=/v0/management/latest-version
Environment=APIM_UPDATE_APPLY_COMMAND=
Environment=APIM_AUTO_UPDATE_ENABLED=true
Environment=APIM_AUTO_UPDATE_SERVICE=cliproxy
Environment=APIM_AUTO_UPDATE_COMPOSE_FILE=
Environment=APIM_AUTO_UPDATE_WORKING_DIR=
Environment=TELEGRAM_POLL_INTERVAL=3s
Environment=HTTP_PROXY=http://127.0.0.1:20171
Environment=HTTPS_PROXY=http://127.0.0.1:20171
Environment=ALL_PROXY=
Environment=NO_PROXY=127.0.0.1,localhost
EnvironmentFile=-/etc/default/cliproxy-access-manager

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/cliproxy-access-manager
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
